# supervisor-bot.conf
[program:Live bot]
directory=/root/bot
command=python3 bot.py
user=root
autostart=true
autorestart=false
stder_logfile=/root/trading-bot.log
stdout_logfile=/root/trading-bot.log